# supervisord.conf
[supervisord]
nodaemon=true

[program:start]
command=/start.sh
